# infra/systemd/juli-secrets-refresh.timer
# Run juli-secrets-refresh.service every 24 hours.
# Install to /etc/systemd/system/juli-secrets-refresh.timer on the production VPS.
[Unit]
Description=Daily Juli production secret refresh

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m
Unit=juli-secrets-refresh.service

[Install]
WantedBy=timers.target
