# infra/systemd/juli-secrets-refresh.service
# Periodic secret synchronization — compare-and-restart only when env files change.
# Install to /etc/systemd/system/juli-secrets-refresh.service on the production VPS.
# Enable the companion timer: systemctl enable --now juli-secrets-refresh.timer
[Unit]
Description=Juli production secret refresh (AWS Secrets Manager)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=AWS_CONFIG_FILE=/etc/aws/config
Environment=AWS_PROFILE=juli-vps-secrets-reader
Environment=AWS_REGION=us-east-2
ExecStart=/root/Juli-AI-v2/infra/scripts/refresh-secrets.sh

[Install]
WantedBy=multi-user.target
